# kavach-07.service
[Unit]
Description=KAVACH-07 Crypto Futures Signal Bot
After=network.target

[Service]
Type=simple
# FIX: Changed from ubuntu to kavach to match install.sh
User=kavach
Group=kavach
WorkingDirectory=/opt/kavach-07
Environment=PYTHONPATH=/opt/kavach-07
EnvironmentFile=/opt/kavach-07/.env
ExecStart=/opt/kavach-07/venv/bin/python -m main
Restart=always
RestartSec=10
StandardOutput=append:/var/log/kavach-07.log
StandardError=append:/var/log/kavach-07.log

[Install]
WantedBy=multi-user.target
